# clients/statrep.service
[Unit]
Description=FAI admins systemd status report

[Service]
Type=oneshot
ExecStart=/usr/local/bin/statrep
